# home-telegram-bot.service
[Unit]
Description=Run Home Telegram Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment="NODE_ENV=production"
ExecStart=/usr/bin/node index.js
WorkingDirectory=
User=
Group=

[Install]
RequiredBy=multi-user.target
